# setup/weather-display.timer
[Unit]
Description=Run weather display update every 2 hours

[Timer]
OnBootSec=1min
OnUnitActiveSec=2h
Persistent=true

[Install]
WantedBy=timers.target
